# deploy/launchd/com.magi.betelgeuse-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Placeholders __PYTHON__ / __ROOT__ / __ENV__ are filled in by deploy/setup-mini.sh
     (__ENV__ from MAGI_ENV in deploy/config.sh, default prod).
     The betelgeuse function's background scheduler. Runs as its OWN process with
     WorkingDirectory inside the function package, so its imports/paths resolve
     exactly as standalone (the isolation contract). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.magi.betelgeuse-worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__ROOT__/functions/betelgeuse/worker.py</string>
        <string>--env</string>
        <string>__ENV__</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__ROOT__/functions/betelgeuse</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__ROOT__/functions/betelgeuse/data/logs/worker.out.log</string>
    <key>StandardErrorPath</key>
    <string>__ROOT__/functions/betelgeuse/data/logs/worker.err.log</string>
</dict>
</plist>
