# deploy/launchd/com.magi.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Placeholders __PYTHON__ / __ROOT__ / __ENV__ are filled in by deploy/setup-mini.sh
     (__ENV__ from MAGI_ENV in deploy/config.sh, default prod).
     The APP-WIDE magi worker: runs host-native functions' scheduled jobs (the Notifier's
     personal reminders today). WorkingDirectory is the magi root so `from functions.* import`
     resolves. MAGI_DATA_DIR points the function send-helpers at the host DB + per-env scope DB
     (app-wide Telegram creds + enable gates). Betelgeuse has its OWN worker (com.magi.betelgeuse-worker). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.magi.worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__ROOT__/worker.py</string>
        <string>--env</string>
        <string>__ENV__</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__ROOT__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>MAGI_ENV</key>
        <string>__ENV__</string>
        <key>MAGI_DATA_DIR</key>
        <string>__ROOT__/data</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__ROOT__/data/logs/worker.out.log</string>
    <key>StandardErrorPath</key>
    <string>__ROOT__/data/logs/worker.err.log</string>
</dict>
</plist>
